# clear-spotify-cache.service
[Unit]
Description=Clear Spotify Cache

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'rm -rf %h/.var/app/com.spotify.Client/cache/spotify/Data/*'

